# percona/percona.ini
[program:percona]
command=/usr/bin/mysqld --pid-file=/run/mysqld/mysqld.pid --bind-address=0.0.0.0 --skip-name-resolve --datadir=/var/lib/mysql
process_name=%(program_name)s
numprocs=1
directory=/var/lib/mysql
autostart=true
autorestart=true
startsecs=10
startretries=3
user=mysql
stdout_logfile=/var/log/supervisor/percona-stdout.log
stderr_logfile=/var/log/supervisor/percona-stderr.log
