# nginx/nginx.ini
[program:nginx]
command=/usr/bin/nginx -g 'pid /run/nginx.pid; daemon off; master_process on;'
process_name=%(program_name)s
numprocs=1
directory=/srv/http
autostart=true
autorestart=true
startsecs=10
startretries=3
user=root
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
