# nginx-hhvm/hhvm.ini
[program:hhvm]
command=/usr/bin/hhvm --user http --mode server -vServer.Type=fastcgi -vServer.Port=9000
process_name=%(program_name)s
numprocs=1
directory=/srv/http
autostart=true
autorestart=true
startsecs=10
startretries=3
user=root
stdout_logfile=/var/log/supervisor/hhvm-stdout.log
stderr_logfile=/var/log/supervisor/hhvm-stderr.log
