# groovebasin/groovebasin.ini
[program:groovebasin]
command=/usr/bin/groovebasin
process_name=%(program_name)s
numprocs=1
directory=/home/groovebasin
autostart=true
autorestart=true
startsecs=10
startretries=3
user=groovebasin
stdout_logfile=/var/log/supervisor/groovebasin-stdout.log
stderr_logfile=/var/log/supervisor/groovebasin-stderr.log
