# nginx-php/php-fpm.ini
[program:php-fpm]
command=/usr/bin/php-fpm --nodaemonize --pid /run/php-fpm/php-fpm.pid
process_name=%(program_name)s
numprocs=1
directory=/srv/http
autostart=true
autorestart=true
startsecs=10
startretries=3
user=root
stdout_logfile=/var/log/supervisor/php-fpm-stdout.log
stderr_logfile=/var/log/supervisor/php-fpm-stderr.log
